# .config/systemd/user/scream.service
[Unit]
Description=Scream service

[Service]
ExecStart=/usr/local/bin/scream

[Install]
WantedBy=default.target
