# .config/systemd/user/swayidle-lock.service
[Unit]
Description=Swayidle lock daemon
BindsTo=sway-session.target

[Service]
ExecStart=/usr/bin/env swayidle -w \
    timeout 660 "loginctl lock-session" \
    before-sleep "loginctl lock-session" \
    lock "/usr/bin/ssh-add -D 2> /dev/null; swaylock"
Restart=always

[Install]
WantedBy=sway-session.target
