# .config/systemd/user/swayidle-screen-off.service
[Unit]
Description=Swayidle screen off daemon
BindsTo=sway-session.target

[Service]
ExecStart=/usr/bin/env swayidle \
    timeout 600 'swaymsg "output * dpms off"' \
    resume 'swaymsg "output * dpms on"'
Restart=always

[Install]
WantedBy=sway-session.target
